# Debian-Ubuntu/vuln-scan.timer
[Unit]
Description=Weekly Vulnerability Scan Timer

[Timer]
OnCalendar=Sun 23:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
